# com.indyhall.sonos-ui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.indyhall.sonos-ui</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/indyhall/.bun/bin/bun</string>
    <string>run</string>
    <string>/Users/indyhall/Sites/sonos-ui/server/index.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/indyhall/Sites/sonos-ui</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/indyhall/logs/sonos-ui.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/indyhall/logs/sonos-ui.error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/indyhall</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
